# bitcoin_supervisor.conf
[program:bitcoind]
command=/snap/bin/bitcoin-core.daemon
autostart=true
autorestart=true
user=vagrant
stderr_logfile=/var/log/bitcoind.err.log
stdout_logfile=/var/log/bitcoind.out.log
